# cron/RasPiCam.cron
# Read data from sensors
*/5 * * * *	/usr/lib/RasPiCamReader.sh
