# cron/RasPiCam_dht11_hack.cron
*/5 * * * * root /usr/lib/RasPiCam_dht11_hack.sh

