# cron/RasPiCam.cron
# Read data from sensors
*/5 * * * * www-data /usr/lib/RasPiCam/RasPiCamReader.sh
